# system_services/hailo-florence/hailo-florence.service
[Unit]
Description=Hailo Florence-2 Captioning + VQA Service
Documentation=https://github.com/hailo-ai/hailo-rpi5-examples
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=hailo-florence
Group=hailo-florence

WorkingDirectory=/var/lib/hailo-florence
StateDirectory=hailo-florence

# XDG wiring: force config + data locations
Environment=XDG_CONFIG_HOME=/etc/xdg
Environment=XDG_CONFIG_DIRS=/etc/xdg
Environment=XDG_DATA_HOME=/var/lib
Environment=XDG_DATA_DIRS=/var/lib:/usr/share:/usr/local/share

Environment=HAILO_PRINT_TO_SYSLOG=1
Environment=PYTHONUNBUFFERED=1

# Service execution
ExecStart=/opt/hailo-florence/venv/bin/python3 /opt/hailo-florence/hailo_florence_service.py

# Resource limits
MemoryMax=3G
CPUQuota=70%
TasksMax=50

# Restart policy
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hailo-florence
SupplementaryGroups=video

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hailo-florence

[Install]
WantedBy=multi-user.target
