# system_services/hailo-scrfd/hailo-scrfd.service
[Unit]
Description=Hailo SCRFD (Face Detection with Landmarks on Hailo-10H)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=hailo-scrfd
Group=hailo-scrfd

# Set working directory for log file creation
WorkingDirectory=/var/lib/hailo-scrfd

# Ensure writable state dir exists and is owned correctly
StateDirectory=hailo-scrfd

# XDG wiring: force config + data locations
Environment=XDG_CONFIG_HOME=/etc/xdg
Environment=XDG_CONFIG_DIRS=/etc/xdg
Environment=XDG_DATA_HOME=/var/lib
Environment=XDG_DATA_DIRS=/var/lib:/usr/share:/usr/local/share

# HailoRT and Python optimization
Environment=HAILO_PRINT_TO_SYSLOG=1
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 /opt/hailo-scrfd/hailo_scrfd_service.py

Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Model loading overhead; extended startup timeout
TimeoutStartSec=120

# Resource limits for face detection
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
