# system_services/hailo-ollama/hailo-ollama.service
[Unit]
Description=Hailo Ollama (Ollama-compatible API on Hailo-10H)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=hailo-ollama
Group=hailo-ollama

# Set working directory for log file creation
WorkingDirectory=/var/lib/hailo-ollama

# Ensure writable state dir exists and is owned correctly
StateDirectory=hailo-ollama

# Bind package manifests into the writable data directory
BindReadOnlyPaths=/usr/share/hailo-ollama/models/manifests:/var/lib/hailo-ollama/models/manifests

# XDG wiring: force config + data locations
Environment=XDG_CONFIG_HOME=/etc/xdg
Environment=XDG_CONFIG_DIRS=/etc/xdg
# XDG_DATA_HOME: writable location for downloaded model blobs
Environment=XDG_DATA_HOME=/var/lib
# XDG_DATA_DIRS: must include /usr/share for package-installed manifests
Environment=XDG_DATA_DIRS=/var/lib:/usr/share:/usr/local/share

# HailoRT logging to journald instead of local files
Environment=HAILO_PRINT_TO_SYSLOG=1
# Enable monitoring for hailortcli monitor
Environment=HAILO_MONITOR=1
Environment=HAILO_TRACE=scheduler

ExecStart=/usr/bin/env hailo-ollama

Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Optional Pi guardrails (tune as needed)
MemoryMax=4G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
