# system_services/hailo-ocr/hailo-ocr.service
[Unit]
Description=Hailo OCR (PaddleOCR on Hailo-10H)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=hailo-ocr
Group=hailo-ocr

# Set working directory for log file creation
WorkingDirectory=/var/lib/hailo-ocr

# Ensure writable state dir exists and is owned correctly
StateDirectory=hailo-ocr

# XDG wiring: force config + data locations
Environment=XDG_CONFIG_HOME=/etc/xdg
Environment=XDG_CONFIG_DIRS=/etc/xdg
Environment=XDG_DATA_HOME=/var/lib
Environment=XDG_DATA_DIRS=/var/lib:/usr/share:/usr/local/share

# Python unbuffered output for real-time logging
Environment=PYTHONUNBUFFERED=1
Environment=OCR_CONFIG_PATH=/etc/xdg/hailo-ocr/hailo-ocr.json

ExecStart=/opt/hailo-ocr/venv/bin/python3 /opt/hailo-ocr/hailo_ocr_server.py

Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Resource limits (tunable via systemd drop-in overrides)
MemoryMax=2.5G
CPUQuota=75%

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hailo-ocr /etc/xdg/hailo-ocr

[Install]
WantedBy=multi-user.target
