# system_services/hailo-whisper/hailo-whisper.service
[Unit]
Description=Hailo Whisper (Speech-to-Text on Hailo-10H)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=hailo-whisper
Group=hailo-whisper

WorkingDirectory=/var/lib/hailo-whisper
StateDirectory=hailo-whisper

Environment=XDG_CONFIG_HOME=/etc/xdg
Environment=XDG_CONFIG_DIRS=/etc/xdg
Environment=XDG_DATA_HOME=/var/lib
Environment=XDG_DATA_DIRS=/var/lib:/usr/share:/usr/local/share

Environment=HAILO_PRINT_TO_SYSLOG=1
Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/hailo-whisper/venv/bin/python3 /opt/hailo-whisper/hailo_whisper_server.py

Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

MemoryMax=3G
CPUQuota=80%

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hailo-whisper /etc/xdg/hailo-whisper

[Install]
WantedBy=multi-user.target
